# infra/mail-processor.service
[Unit]
Description=MCP Mail Processor (Unified System)
After=network-online.target
Wants=network-online.target

[Service]
User=gonya
WorkingDirectory=/home/gonya/Documents/Unified_System_Core
ExecStart=/run/current-system/sw/bin/bash -lc 'cd /home/gonya/Documents/Unified_System_Core && /run/current-system/sw/bin/devenv shell -- python3 Scripts/Orchestration/mail_processor.py'
Restart=always
RestartSec=10
EnvironmentFile=/home/gonya/Documents/Unified_System_Core/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
