# Scripts/automation/alice-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Alice Skill
After=network.target ai-bot.service
Wants=ai-bot.service

[Service]
Type=simple
User=gonya
WorkingDirectory=/home/gonya/Documents/Unified_System
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:8090 --no-autoupdate
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="HOME=/home/gonya"

[Install]
WantedBy=multi-user.target
