# hooks/com.antigravity.hooks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.hooks</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/macbook/Documents/Unified_System/hooks/hooks_manager.py</string>
        <string>--serve</string>
        <string>--port</string>
        <string>8766</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/macbook/Documents/Unified_System/hooks/hooks_server.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbook/Documents/Unified_System/hooks/hooks_server.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/macbook/Documents/Unified_System/hooks</string>
</dict>
</plist>
