# Projects/AI_Core/scripts/ai-bot.service
[Unit]
Description=AI Telegram Bot (Unified System)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Unified_System
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=%h/Unified_System/Projects/AI_Core/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/Unified_System/Projects/AI_Core/venv/bin/python3 Projects/AI_Core/src/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
