# Scripts/Orchestration/com.unified.mail-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.unified.mail-processor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/macbook/Documents/Unified_System/.venv/bin/python</string>
        <string>/Users/macbook/Documents/Unified_System/Scripts/Orchestration/mail_processor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/macbook/Documents/Unified_System/Scripts/Orchestration</string>
    <key>StandardOutPath</key>
    <string>/Users/macbook/Documents/Unified_System/Reports/mail_processor.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbook/Documents/Unified_System/Reports/mail_processor.err</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>AGENT_MAIL_NAME</key>
        <string>PinkLake</string>
        <key>AGENT_MAIL_PROJECT</key>
        <string>/home/gonya/Unified_System</string>
        <key>AGENT_MAIL_SERVER</key>
        <string>http://100.110.209.49:8765</string>
    </dict>
</dict>
</plist>
